# db/container/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:postgres1c]
command=/runpg.sh
user=postgres
startsecs=5
autostart=true