# server/container/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:server1c]
command="/run1c.sh"
startsecs=5
autostart=true